# build/tibero/src/supervisord.conf
[supervisord]
nodaemon=true
logfile = $TB_MOUNT_VOLUME_PATH/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info

[program:run]
startretries=3
command=/bin/bash -c "$SRC_HOME/boot.sh > $TB_MOUNT_VOLUME_PATH/cmdlog 2>&1"
autorestart=false
stopasgroup=true
user=hyperdata
stderr_logfile=$TB_MOUNT_VOLUME_PATH/super_stderr.log
stdout_logfile=$TB_MOUNT_VOLUME_PATH/super_stdout.log
