# build/tibero/src/supervisord.conf
[supervisord]
nodaemon=true
logfile = /db/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info

[program:run]
command=/bin/bash -c "/deploy_src/src/boot.sh > /db/cmdlog 2>&1"
autorestart=false
stopasgroup=true
user=hyperdata
stderr_logfile=/db/super_stderr.log
stdout_logfile=/db/super_stdout.log
